# packaging/org.tizen.fido.service
[Unit]
Description=FIDO Service D-Bus

[D-BUS Service]
Name=org.tizen.fido
Exec=/usr/bin/fido-service
User=system
Group=system
