# packaging/org.tizen.dummyasm.service
[Unit]
Description=FIDO Dummy ASM D-Bus

[D-BUS Service]
Name=org.tizen.dummyasm
Exec=/usr/bin/dummyasm-service
User=service_fw
Group=service_fw
